# vsftpd@.service
[Unit]
Description=Vsftpd ftp daemon
After=network-online.target
PartOf=vsftpd.target

[Service]
Type=forking
ExecStartPre=/usr/libexec/vsftpd-ssl-gencerts
ExecStart=/usr/sbin/vsftpd /etc/vsftpd/%i.conf

[Install]
WantedBy=vsftpd.target
